# lib/systemd/system/tor@default.service.d/30_clean-torrc-d-on-reload.conf
## Copyright (C) 2019 - 2019 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

#### meta start
#### project Whonix
#### category tor
#### gateway_only yes
#### description
## Run torrc.d cleaner at reload.
## https://forums.whonix.org/t/torrc-d-is-comming/4041/72
[Service]
ExecReload=
ExecReload=sh -c "/usr/lib/anon-gw-anonymizer-config/torrc-d-cleaner && /bin/kill -HUP ${MAINPID}"

#### meta end
