# lib/systemd/system/tor@default.service.d/40_network-target.conf
## Copyright (C) 2023 - 2023 ENCRYPTED SUPPORT LP <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

#### meta start
#### project Whonix
#### category tor
#### gateway_only yes
#### description

## Workaround for upstream bug.
## https://forums.whonix.org/t/failed-to-start-anonymizing-overlay-network-for-tcp-tor-fails-to-start-a-few-times-before-succeeding-to-start/16289/12
## https://bugs.debian.org/cgi-bin/bugreport.cgi?bug=1029554
[Unit]
After=network.target

#### meta end
