# lib/systemd/system/tor@default.service.d/40_obfs4proxy-workaround.conf
## The hardening option NoNewPrvileges=Yes in the Debian tor@default.service file breaks loading of obfs4proxy in combination with Debian's apparmor profile which uses PUx. (ix execution permissions work with NoNewPrivileges=No)
[Service]
NoNewPrivileges=
