# lib/systemd/system/tor@default.service.d/40_obfs4proxy-workaround.conf
## Copyright (C) 2018 - 2019 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

#### meta start
#### project Whonix
#### category tor
#### gateway_only yes
#### description

## The hardening option NoNewPrvileges=Yes in the Debian tor@default.service file breaks loading of obfs4proxy in combination with Debian's apparmor profile which uses PUx. (ix execution permissions work with NoNewPrivileges=No)
## https://phabricator.whonix.org/T676
[Service]
NoNewPrivileges=no

#### meta end
