# lib/systemd/system/tor@default.service.d/40_torrcd-workaround.conf
## Tor will fail to start if a specified torrc file is missing. There is no way to let Tor ignore the missing torrc files if they are specified by %include approach. Therefore, we have to make sure all the torrc files exist before starting Tor with systemd.
##
[Unit]
After=anon-gw-anonymizer-config.service
Wants=anon-gw-anonymizer-config.service
