# lib/systemd/system/tor@default.service.d/50_controlsocket-workaround.conf
## When both User option and DisabkeNetwork 1 are set as part of the Tor configurations, the Tor cotrol socket will not be opened unless reloading Tor.
## https://lists.torproject.org/pipermail/tor-dev/2018-February/012887.html
[Service]
ExecStartPost=/bin/kill -HUP ${MAINPID}