# lib/systemd/system/vanguards.service.d/30_anon-gw-anonymizer-config.conf
## Copyright (C) 2019 - 2019 ENCRYPTED SUPPORT LP <adrelanos@riseup.net>
## See the file COPYING for copying conditions.

[Service]
Environment=VANGUARDS_CONFIG=/etc/tor/anon-vanguards.conf
